# acinstaller
#!/sbin/openrc-run

start() {

	dialog  --colors \
		--backtitle "\Z1\Zu\ZbAC-Core Installer\Zn" \
		--title "Welcome To Your AC-Core Node!" \
		--msgbox 'Lets get your node setup!' \
		10 40;

	dialog  --colors \
		--backtitle "\Z1\Zu\ZbAC-Core Installer\Zn" \
		--menu \
		"Select the components to install:" \
		10 40 3 \
		"1" "AC Core"
}
